# catkin_ws/devel/etc/catkin/profile.d/15.rocon_launch.zsh
# generated from rocon_launch/env-hooks/15.rocon_launch.zsh.em

. "/home/giovani/catkin_ws/src/rocon_tools/rocon_launch/shells/env.zsh"
