# catkin_ws/devel/etc/catkin/profile.d/15.rocon_test.zsh
# generated from rocon_test/env-hooks/15.rocon_test.zsh.em

. "/home/giovani/catkin_ws/src/rocon_multimaster/rocon_test/shells/env.zsh"
